# dot_config/systemd/user/chezmoi-update.service
[Unit]
Description=Update dotfiles

[Service]
ExecStart=/home/linuxbrew/.linuxbrew/bin/chezmoi update
